# prefect/supervisord/supervisord.conf
[supervisord]
nodaemon=true

[program:prefect]
command=prefect server start --host 0.0.0.0
